# ha-addon/fitness-habit-tracker/services.d/nginx/run
#!/usr/bin/with-contenv sh
set -eu
echo "==> [nginx] starting"
exec nginx -g "daemon off;"
