# ha-addon/fitness-habit-tracker/services.d/app/run
#!/usr/bin/with-contenv sh
set -eu

# HA add-on persistent storage
mkdir -p /data

# Default DB if not provided
export DATABASE_URL="${DATABASE_URL:-file:/data/dev.db}"
export NODE_ENV="${NODE_ENV:-production}"

echo "==> [app] NODE_ENV=${NODE_ENV}"
echo "==> [app] DATABASE_URL=${DATABASE_URL}"

cd /app

echo "==> [app] prisma generate"
npx prisma generate

echo "==> [app] prisma migrate deploy"
npx prisma migrate deploy

exec npm run start
